# deploy/com.bullbear.paper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bullbear.paper</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Volumes/MOVESPEED/Documents/GitHub/bullbear-ai-trader/.venv/bin/python</string>
    <string>-m</string>
    <string>src.cli</string>
    <string>run-paper</string>
    <string>--config</string>
    <string>config/default.yaml</string>
    <string>--agent</string>
    <string>external</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Volumes/MOVESPEED/Documents/GitHub/bullbear-ai-trader</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/Volumes/MOVESPEED/Documents/GitHub/bullbear-ai-trader/reports/runtime/launchd.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Volumes/MOVESPEED/Documents/GitHub/bullbear-ai-trader/reports/runtime/launchd.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
